# flask_test/wireguard-test.service
[Unit]
Description=WireGuard Test Flask Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/flask_test
ExecStart=/home/ubuntu/flask_test/venv/bin/python app.py
Restart=always
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
